# deploy/ec2/hydraflow.service
[Unit]
Description=HydraFlow orchestrator and dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hydraflow
Group=hydraflow
WorkingDirectory=/opt/hydraflow
EnvironmentFile=-/etc/hydraflow.env
Environment=HYDRAFLOW_HOME_DIR=/var/lib/hydraflow
RuntimeDirectory=hydraflow
ExecStart=/opt/hydraflow/deploy/ec2/deploy-hydraflow.sh run
Restart=always
RestartSec=5
TimeoutStopSec=60
KillSignal=SIGINT
StandardOutput=append:/var/log/hydraflow/orchestrator.log
StandardError=append:/var/log/hydraflow/orchestrator.log
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
