# systemd/slack-bot.service
[Unit]
Description=Insight Mesh Slack Bot
Documentation=https://github.com/your-org/ai-slack-bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=slack-bot
Group=slack-bot
WorkingDirectory=/opt/slack-bot
ExecStart=/opt/slack-bot/venv/bin/python app.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutSec=30
Restart=always
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=ENVIRONMENT=production
Environment=LOG_LEVEL=INFO
EnvironmentFile=/opt/slack-bot/.env

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/slack-bot/logs
ReadOnlyPaths=/opt/slack-bot

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slack-bot

[Install]
WantedBy=multi-user.target